# systemd/leddy.service
[Unit]
Description=Leddy

[Service]
ExecStart=/home/pi/programming/leddy/Leddy

[Install]
WantedBy=multi-user.target
